# Semester4/B-DOP-400/octopus/roles/result/files/result.service
[Unit]
Description=Result

[Service]
User=admin
WorkingDirectory=/opt/result
ExecStart=sudo /opt/result/run.sh
Restart=always
EnvironmentFile="/opt/result/.env.result"
